# systemd/mark-core-v2.service
[Unit]
Description=Mark Core v2 Backend Service
Documentation=https://github.com/CranckThatFranck/markv2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jarvis
Group=jarvis
WorkingDirectory=/opt/jarvis/backend
ExecStart=/opt/jarvis/backend/.venv/bin/uvicorn src.backend.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
UMask=0027
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mark-core-v2

Environment="MARK_STATE_DIR=/var/lib/jarvis-mark/estado"
Environment="MARK_LOG_DIR=/var/log/jarvis"
EnvironmentFile=-/etc/mark-core-v2/environment
StateDirectory=jarvis-mark
LogsDirectory=jarvis

# Security hardening
PrivateTmp=yes
ProtectSystem=yes
ProtectHome=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
