# packaging/deb/postinst
#!/bin/bash
set -e

# Debian postinst script for mark-core-v2
# This script runs after the package is installed

case "$1" in
    configure)
        # Create jarvis user if it doesn't exist
        if ! id -u jarvis > /dev/null 2>&1; then
            useradd --system --home /var/lib/jarvis-mark --no-create-home jarvis
        fi

        # Create necessary directories
        mkdir -p /var/lib/jarvis-mark
        mkdir -p /var/lib/jarvis-mark/estado
        mkdir -p /var/log/jarvis
        mkdir -p /etc/mark-core-v2

        # Set ownership
        chown -R jarvis:jarvis /opt/jarvis
        chown -R jarvis:jarvis /var/lib/jarvis-mark
        chown -R jarvis:jarvis /var/log/jarvis

        # Set permissions
        chmod 755 /var/lib/jarvis-mark
        chmod 755 /var/log/jarvis
        chmod 755 /var/lib/jarvis-mark/estado

        # Build the runtime virtualenv in the final installed path so entrypoint
        # scripts such as uvicorn get stable shebangs.
        python3.12 -m venv /opt/jarvis/backend/.venv
        /opt/jarvis/backend/.venv/bin/pip install -U pip setuptools wheel
        /opt/jarvis/backend/.venv/bin/pip install -r /opt/jarvis/backend/requirements-backend.txt

        # Keep the packaged environment file readable only by privileged users.
        chmod 640 /etc/mark-core-v2/environment || true

        systemctl daemon-reload || true

        echo "Mark Core v2 installed successfully!"
        echo "To start the service: systemctl start mark-core-v2"
        echo "To enable on boot: systemctl enable mark-core-v2"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
