# packaging/deb/prerm
#!/bin/bash
set -e

# Debian prerm script for mark-core-v2
# This script runs before the package is removed

case "$1" in
    remove|deconfigure)
        # Stop the service if running
        systemctl stop mark-core-v2 || true
        systemctl disable mark-core-v2 || true
        systemctl daemon-reload || true
        ;;

    upgrade)
        # On upgrade, stop the service gracefully
        systemctl stop mark-core-v2 || true
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
